# docker/crontab/crontab.sh
#!/bin/sh

curl "http://127.0.0.1:20000/api/publish/crontab" >> /dev/null 2>&1
